# debian/mythtv-backend.upstart
# MythTV Backend service

description     "MythTV Backend"
author          "Mario Limonciello <superm1@ubuntu.com>"

start on (local-filesystems and net-device-up IFACE=lo and started udev-finish)
stop on starting shutdown

#expect fork
respawn

exec /usr/bin/mythbackend --logfile /var/log/mythtv/mythbackend.log --user mythtv
