# debian/mythtv-backend.cron.weekly
#!/bin/sh
#Janitor mode to clean up stale stuff
DIRECTORY=$(getent passwd mythtv | cut -d':' -f6)
if [ -f "$DIRECTORY/.mythtv/config.xml" ] && [ -x /usr/bin/mythbackend ]; then
    su mythtv -c "/usr/bin/python /usr/share/doc/mythtv-backend/contrib/metadata/jamu.py -MJ >> '/var/log/mythtv/jamu.log' 2>&1"
fi
