# debian/mythvideo.cron.daily
#!/bin/sh
#Daily maintenance mode for Jamu
DIRECTORY=$(getent passwd mythtv | cut -d':' -f6)
if [ -f "$DIRECTORY/.mythtv/config.xml" ] && [ -x /usr/bin/mythbackend ]; then
    su mythtv -c "/usr/bin/python /usr/share/mythtv/mythvideo/scripts/jamu.py -M >> '/var/log/mythtv/jamu.log' 2>&1"
fi

